# systemd/isponsorblocktv-webui.service
[Unit]
Description=iSponsorBlockTV WebUI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Replace REPLACE_ME with the user that owns the pipx install (e.g. dietpi, pi, root).
User=REPLACE_ME
# Path to the pipx-installed shim. pipx puts it under that user's home dir.
# For root, change this to /root/.local/bin/isponsorblocktv-webui.
ExecStart=/home/REPLACE_ME/.local/bin/isponsorblocktv-webui
Restart=on-failure
RestartSec=3
Environment=WEBUI_HOST=0.0.0.0
Environment=WEBUI_PORT=8099
# Uncomment and edit if iSponsorBlockTV runs under a different unit/container name:
# Environment=WEBUI_SERVICE_NAME=iSponsorBlockTV

[Install]
WantedBy=multi-user.target
